# data/backups/service_update_20250818_203712/data/backups/service_update_20250818_203712/kata-web.service
[Unit]
Description=Kata Assistant Web Server
After=network.target

[Service]
Type=simple
User=steveen
WorkingDirectory=/home/steveen/asistente_kata
ExecStart=/home/steveen/gcp-tts-venv/bin/python3 web_server.py
Restart=always

[Install]
WantedBy=multi-user.target
